# s6-overlay/s6-rc.d/fix-chrome-sandbox/run
#!/usr/bin/with-contenv bash
set -euo pipefail

if [ -e /opt/obsidian/chrome-sandbox ]; then
  chown root:root /opt/obsidian/chrome-sandbox
  chmod 4755 /opt/obsidian/chrome-sandbox
fi
